# deploy/kiron-api.service
# systemd unit for the Kiron FastAPI backend.
# Install:  sudo cp deploy/kiron-api.service /etc/systemd/system/
#           sudo systemctl daemon-reload && sudo systemctl enable --now kiron-api
#
# Assumes the repo is cloned at /opt/kiron and a venv exists at
# /opt/kiron/backend/.venv (see DEPLOY.md). Adjust paths if you deploy elsewhere.

[Unit]
Description=Kiron API (FastAPI/uvicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=exec
User=kiron
Group=kiron
WorkingDirectory=/opt/kiron/backend
EnvironmentFile=/etc/kiron/backend.env
# Bind to loopback only — nginx is the public entrypoint and proxies /api here.
# --workers 1: the WebSocket hub (app/routers/ws.py) holds an in-memory dict
# of {user_id -> set[WebSocket]}, so cross-worker broadcasts silently miss
# (sender's POST hits worker A, recipient's socket on worker B never hears).
# At current scale (~30 active users) one worker is plenty. If concurrent
# load ever justifies horizontal scaling, migrate the hub to Redis pub/sub
# per the comment in ws.py before bumping workers back up.
ExecStart=/opt/kiron/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8787 --workers 1
Restart=on-failure
RestartSec=3
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
